# dashboard_api.service
[Unit]
Description=Alpha-Z Trading Dashboard API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vivik
ExecStart=/usr/bin/python3 /root/vivik/dashboard_api.py
Restart=always
RestartSec=10
StandardOutput=append:/root/vivik/logs/dashboard_api.log
StandardError=append:/root/vivik/logs/dashboard_api_error.log

# 환경 변수
Environment="PYTHONUNBUFFERED=1"

# 리소스 제한
MemoryLimit=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
